# units/bootstrap.service
[Unit]
Description=Node bootstrapping

[Service]
Type=oneshot
User=cluster
ExecStart=+/bin/bash /var/lib/cluster/bootstrap.sh
