# units/bootstrap.service
[Unit]
Description=Node bootstrapping

[Service]
Type=oneshot
ExecStart=/bin/bash /local/bootstrap.sh
