# units/bootstrap.service
[Unit]
Description=Node bootstrapping

[Service]
Type=oneshot
ExecStart=/bin/bash /var/lib/cluster/bootstrap.sh
